# deploy/systemd/esg-score-api.service
[Unit]
Description=ESG Score API
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/esg-score/backend
EnvironmentFile=/opt/esg-score/backend/.env
ExecStart=/opt/esg-score/backend/api
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
